# Python_C++_caffe-tfskills-android/cmake-gtest-pytest-ninja-android/ros-cross-compile/ros/my-task-code/20160712-PC/20160712-PC/20160712-PC/startup-check/startup/lu-ban-movebase.conf
author "Haijun Zhao "
description "lu-ban upstart script"
start on stopped rc RUNLEVEL=[2345] and startup and start_movebase

pre-start script
end script

script
cd /home/linaro/nav/src/startup
./start_movebase.bash   2>&1  & 
end script




